# recipes-inobram/inobram-installer/files/inobram-installer.service
[Unit]
Description=Inobram Installer
After=local-fs.target

[Service]
ExecStart=/usr/bin/inobram-installer

[Install]
WantedBy=multi-user.target
