# recipes-smaai5/smaai5-installer/smaai5-installer/smaai5-installer.service
[Unit]
Description=SMAAI5 Installer
After=local-fs.target

[Service]
ExecStart=/usr/bin/smaai5-installer

[Install]
WantedBy=multi-user.target
