# recipes-kernel/bb-cape-overlays/bb-cape-overlays/beaglebone-universal-io.service
[Unit]
Description=Beaglebone Universal IO

[Service]
ExecStart=/usr/bin/config-pin -f /etc/config-pin.conf

[Install]
WantedBy=multi-user.target
